# systemd/user/change_theme.service
[Unit]
Description= changes theme
    
[Service]
Type=oneshot
ExecStart=/home/user/.config/script/theme.sh