# systemd/user/start_hyprland.service
[Unit]
Description= changes theme
    
[Service]
Type=oneshot
ExecStart=hyprland

[Install]
WantedBy=default.target